# configs/crontab/crontab.config
@reboot ( bash ~/node-cms/scripts/production/mongod.sh start && npm start ) >> ~/node-cms/logs/cron.log 2>&1
0 3 * * * bash ~/node-cms/configs/crontab/crontab.sh >> ~/node-cms/logs/cron.log 2>&1
