# configs/crontab/crontab.config
0 3 * * * bash ~/fully-hatter/configs/crontab/crontab.sh >> ~/fully-hatter/log/cron.log 2>&1
