# config/crontab/crontab.config
0 3 * * * bash /home/ubuntu/fully-hatter/config/crontab/crontab.sh >> /home/ubuntu/fully-hatter/log/cron.log 2>&1
